# roles/superfans/templates/superfans.service.j2
[Unit]
Description=Superfans

[Service]
EnvironmentFile=-{{ superfans_env_file }}
ExecStart={{ superfans_wrapper }}

[Install]
WantedBy=multi.user.target
